# workers.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
username=trans            ; (default is no username (open server))
password=e20ecac90375f3   ; (default is no password (open server))

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9003        ; (ip_address:port specifier, *:port for all iface)
username=trans             ; (default is no username (open server))
password=e20ecac90375f3            ; (default is no password (open server))

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
umask=022                   ; (process file creation umask;default 022)
user=root                ; (default is current user, required if root)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
serverurl=http://127.0.0.1:9003 ; use an http:// url to specify an inet socket
username=trans             ; should be same as http_username if set
password=e20ecac90375f3               ; should be same as http_password if set
prompt=trans        ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available


[program:stats]
process_name=%(process_num)s
command=php /home/dev/current/artisan queue:work --queue=stats --sleep=1 --timeout=900
autostart=true
autorestart=true
user=www-data
numprocs=6
;redirect_stderr=true
;stdout_logfile=/home/dev/current/storage/logs/queue-stats.log

[program:listings]
process_name=%(process_num)s
command=php /home/dev/current/artisan queue:work --queue=listings --sleep=3 --timeout=900
autostart=true
autorestart=true
user=www-data
numprocs=4
;redirect_stderr=true
;stdout_logfile=/home/dev/current/storage/logs/queue-listings.log
